# upstart/waggle_plugin_manager.conf
start on (local-filesystems)
stop on runlevel[016]
respawn
respawn limit unlimited
post-stop exec sleep 5
chdir /usr/lib/waggle/plugin_manager
exec python plugin_manager.py --force --logging
